# dot_config/systemd/user/theme-sync.service
[Unit]
Description=Sync dev tool themes with system color scheme
After=plasma-workspace.target
Wants=plasma-workspace.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/theme-sync
Restart=on-failure
RestartSec=5

[Install]
WantedBy=plasma-workspace.target
